# install/package/nas/aaa_02/fstatic/etc/cron.daily/a0_TaskDaily
#!/bin/sh

cd /usr/lib/scripter/backup
./backupme.sh
