# docker/lamp2/docker/dist/python_crawler/python/f_static/etc/systemd/system/vCrawler.service
[Unit]
Description=vCrawler service
After=network.target

[Service]
EnvironmentFile=/etc/default/vCrawler
WorkingDirectory=/usr/lib/py-vCrawler/src
User=admin
#ExecStart=/usr/bin/python3 -B vCrawler.py
ExecStart=/bin/bash -c './vCrawler.sh'
Restart=always
RestartSec=15
Type=simple

[Install]
WantedBy=multi-user.target
