# docker/lamp2/docker/dist/python_ssh_scraper/python/f_static/etc/systemd/system/vScraper.service
[Unit]
Description=vScraper service
After=network.target

[Service]
EnvironmentFile=/etc/default/vScraper
WorkingDirectory=/usr/lib/vScraper/src
User=admin
ExecStart=/usr/bin/python3 -B vScraper.py
Restart=always
RestartSec=15
Type=simple

[Install]
WantedBy=multi-user.target
